# katana-finality-tracker.service
[Unit]
Description=Katana Finality Tracker
Documentation=https://github.com/0xsharma/katana-finality-tracker
After=network.target
Wants=network.target

[Service]
Type=simple
WorkingDirectory=/opt/katana-finality-tracker
ExecStart=/opt/katana-finality-tracker/katana-finality-tracker
Restart=always
RestartSec=10
StandardOutput=journal
StandardError=journal
SyslogIdentifier=katana-finality-tracker

# Environment variables
EnvironmentFile=/opt/katana-finality-tracker/.env

# Resource limits
LimitNOFILE=65536
LimitNPROC=4096

[Install]
WantedBy=multi-user.target 
